# conf/supervisord.conf
[supervisord]
nodaemon=true

[program:tomcat]
command=/opt/liferay/init.sh
numprocs=1
autostart=true
autorestart=true